# deployment_scripts/puppet/modules/calico/templates/calico-alt-gateway.conf.erb
description "calico-alt-gateway"
author "Alternative default gateway for Calico network"

#start on stopped rc RUNLEVEL=[2345]
start on (net-device-up and started etcd)
#stop on runlevel [!2345]

#respawn

script
  ip r add default via <%= @calico_alt_gateway %> table <%= @calico_mark %>
  ip rule add fwmark <%= @calico_mark %> table <%= @calico_mark %>
end script
